# FW/patch/etc/rc.local
#!/bin/sh

echo "initramfs test"

modprobe cmac
modprobe bnep
modprobe 8021q
modprobe vc4
modprobe brcmfmac
modprobe v3d
modprobe raspberrypi_hwmon
modprobe bcm2835_codec
modprobe bcm2835_v4l2
